# simple-service/simple-service.service
[Unit]
Description=Simple service example

[Service]
Type=simple
ExecStart=/usr/bin/simple-service

[Install]
WantedBy=multi-user.target
